# deployment/clarifier-api.service
[Unit]
Description=Argument Clarifier API Server
Documentation=https://github.com/mythofkas-commits/bias-clarity-extension
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/argument-clarifier/server
Environment=NODE_ENV=production
EnvironmentFile=/var/www/argument-clarifier/server/.env
ExecStart=/usr/bin/node /var/www/argument-clarifier/server/dist/index.js

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/argument-clarifier/server

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clarifier-api

[Install]
WantedBy=multi-user.target
